# launchd/com.YOURUSER.update-all-clis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for update-all-clis
  Copy this to ~/Library/LaunchAgents/
  Edit the ProgramArguments path to match your installation location.
  Then load with: launchctl load ~/Library/LaunchAgents/com.YOURUSER.update-all-clis.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.YOURUSER.update-all-clis</string>

	<key>ProgramArguments</key>
	<array>
		<!-- EDIT THIS LINE: set the full path to update_all_clis.sh on your system -->
		<string>/FULL/PATH/TO/update-all-clis/scripts/update_all_clis.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>8</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<!-- Logs stdout/stderr to files -->
	<key>StandardOutPath</key>
	<string>/Users/YOURUSER/.hermes/logs/update-all-clis.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOURUSER/.hermes/logs/update-all-clis.err</string>

	<!-- Don't run on login; only at the scheduled time -->
	<key>RunAtLoad</key>
	<false/>

	<key>KeepAlive</key>
	<false/>
</dict>
</plist>
