# launchd/com.YOURUSER.update-all-clis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for update-all-clis
  Copy this to ~/Library/LaunchAgents/
  Edit the ProgramArguments path to match your installation location.
  Then load with: launchctl load ~/Library/LaunchAgents/com.YOURUSER.update-all-clis.plist

  Logs: ~/.config/update-all-clis/logs/update-all-clis.log (and .err)
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.YOURUSER.update-all-clis</string>

	<key>ProgramArguments</key>
	<array>
		<!-- EDIT: full path to update_all_clis.sh (install.sh defaults to ~/update-all-clis/) -->
		<string>/FULL/PATH/TO/update-all-clis/update_all_clis.sh</string>
	</array>

	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>8</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/YOURUSER/.config/update-all-clis/logs/update-all-clis.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/YOURUSER/.config/update-all-clis/logs/update-all-clis.err</string>

	<key>RunAtLoad</key>
	<false/>

	<key>KeepAlive</key>
	<false/>

	<!-- No GUI dialog for scheduled runs (manual runs use Terminal TTY) -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>UPDATE_ALL_CLIS_NO_NOTIFY</key>
		<string>1</string>
	</dict>
</dict>
</plist>
